# v2/deploy/supervisord.conf
[supervisord]
logfile=/Users/joshcrasto/Projects/lichess_analysis/v2/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/Users/joshcrasto/Projects/lichess_analysis/v2/logs/supervisord.pid
nodaemon=false
user=joshcrasto

[unix_http_server]
file=/Users/joshcrasto/Projects/lichess_analysis/v2/logs/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///Users/joshcrasto/Projects/lichess_analysis/v2/logs/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ── Backend ───────────────────────────────────────────────────────────────────

[program:backend]
command=/Users/joshcrasto/Projects/lichess_analysis/v2/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --log-level info
directory=/Users/joshcrasto/Projects/lichess_analysis/v2/backend
autostart=true
autorestart=true
startretries=3
startsecs=3
stopwaitsecs=10
stdout_logfile=/Users/joshcrasto/Projects/lichess_analysis/v2/logs/backend.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/Users/joshcrasto/Projects/lichess_analysis/v2/logs/backend.err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
environment=PATH="/Users/joshcrasto/Projects/lichess_analysis/v2/backend/.venv/bin:/usr/local/bin:/usr/bin:/bin"

; ── Frontend ─────────────────────────────────────────────────────────────────

[program:frontend]
command=/Users/joshcrasto/.nvm/versions/node/v24.14.0/bin/npm run dev
directory=/Users/joshcrasto/Projects/lichess_analysis/v2/frontend
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
stdout_logfile=/Users/joshcrasto/Projects/lichess_analysis/v2/logs/frontend.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/Users/joshcrasto/Projects/lichess_analysis/v2/logs/frontend.err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
environment=PATH="/Users/joshcrasto/.nvm/versions/node/v24.14.0/bin:/usr/local/bin:/usr/bin:/bin"
